# root/etc/s6-overlay/s6-rc.d/svc-icd/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash


exec s6-setuidgid abc /app/init.sh